# prod/rootfs/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/scripts/entrypoint-nginx.sh
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log