# prod/rootfs/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpmPHP_VERSION_HERE -F
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm.log